# tools/ts-specs/ext-sources/zabbix_server.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='zabbix-server'>
  <service name='network/zabbix/server' type='service' version='0'>

    <create_default_instance enabled="false" />
    <single_instance/>

    <dependency name='multi-user-server' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user-server'/>
    </dependency>

    <exec_method name='start' type='method' exec='/opt/ts/sbin/zabbix_server' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec='pkill zabbix_server' timeout_seconds='60'/>

    <template>
      <common_name>
        <loctext xml:lang='C'>Zabbix server</loctext>
      </common_name>
    </template>

  </service>
</service_bundle>
