# wikitex/wikitex.in.cron
MAILTO=%MAILTO%
# At midnight, clean up WikiTeX cruft.
0 0 * * *      %DECRUFT% 2>&1
