# debs/wikimedia-base/debian/postinst
#! /bin/sh
# postinst script for wikimedia-base
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

. /lib/lsb/init-functions

case "$1" in
    configure)
	# Determine which cluster this host is part of
	if [ -f /etc/wikimedia-cluster ]
	then
		CLUSTERNAME=$(cat /etc/wikimedia-cluster)
	else
		# FIXME: we can (need?) to do better than this
		case $(hostname -d) in
		knams.wikimedia.org|esams.wikimedia.org)
			CLUSTERNAME=esams
			;;
		pmtpa.wmnet|wikimedia.org)
			CLUSTERNAME=pmtpa
			;;
		esac
		
		echo $CLUSTERNAME > /etc/wikimedia-cluster
	fi

	if ! grep -qs "^Explanation: Prefer Wikimedia APT repository packages in all cases" /etc/apt/preferences
	then
		# Set package pinning to prefer Wikimedia APT repository packages
		[ -f /etc/apt/preferences ] && echo >> /etc/apt/preferences
		cat /usr/share/doc/wikimedia-base/apt_preferences >> /etc/apt/preferences
	fi

	# Disable the 'quiet' kernel command line option so console messages
	# will be printed.
	if [ -f /boot/grub/menu.lst ] && grep -q "^# defoptions.*[= ]quiet " /boot/grub/menu.lst
	then
		sed -i '/^# defoptions.*[= ]quiet /s/quiet //' /boot/grub/menu.lst
		update-grub
	elif [ -f /etc/default/grub ] && grep -q '^GRUB_CMDLINE_LINUX_DEFAULT="quiet splash"' /etc/default/grub
	then
		sed -i '/^GRUB_CMDLINE_LINUX_DEFAULT="quiet splash"/s/quiet splash//' /etc/default/grub
		update-grub
	fi

	# Optimize TCP
	if ! grep -q "^### START: Added by wikimedia-base" /etc/sysctl.conf || dpkg --compare-versions "$2" lt "0.18"
	then
		cat /usr/share/doc/wikimedia-base/sysctl.conf >> /etc/sysctl.conf
		sysctl -p /usr/share/doc/wikimedia-base/sysctl.conf
	fi

	# Default editor joe? wtf!?
	update-alternatives --install /usr/bin/editor editor /usr/bin/vim 75 \
		--slave /usr/share/man/man1/editor.1.gz editor.1.gz \
		/usr/share/man/man1/vim.1.gz

	# Enable some aliases in /root/.bashrc
	sed -i '
		/^#alias ll=/ s/^#//
		/^#alias la=/ s/^#//
	' /root/.bashrc

	# Enable system accounting
	sed -i 's/^ENABLED="false"/ENABLED="true"/g' /etc/default/sysstat
	invoke-rc.d sysstat start || true
	
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

