# mwsearch/scripts/mwsearchd.sh.in
#!/bin/sh
# chkconfig: 2345 84 16
# description: MediaWiki search backend HTTP server
# vim: autoindent

SERVICE_NAME="MediaWiki search server"
BINDIR=/usr/local/bin
. /etc/init.d/functions

case "$1" in
	start)
		# Check if it is running
		pid=`pidfileofproc mwsearchd`
		if [ ! -z $pid ] && checkpid $pid; then
			echo "$SERVICE_NAME already running."
			exit 1
		fi

		echo -n "Starting $SERVICE_NAME: "

		# Increase FD limit
		ulimit -n 8192
		# Run the daemon
		# Needs daemonize from http://www.clapper.org/software/daemonize/
		daemonize -p /var/run/mwsearchd.pid $BINDIR/mwsearchd
		[ "$?" -eq 0 ] && success $SERVICE_NAME || failure $SERVICE_NAME
		echo
		;;
	stop)
		echo -n "Stopping $SERVICE_NAME: "
		killproc mwsearchd
		echo
		;;
	restart)
		$0 stop
		$0 start
		;;
	status)
		pid=`pidfileofproc mwsearchd`
		if [ ! -z $pid ] && checkpid $pid; then
			echo "$SERVICE_NAME is running"
		else
			echo "$SERVICE_NAME is not running"
		fi
		;;
	*)
		echo "Usage: $0 {start|stop|status|restart}"
		exit 1
		;;
esac

