# tools/invoice-generator/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 0 10 * * php /home/tstarling/src/mediawiki/trunk/tools/invoice-generator/invoice.php
