# tools/smf/fisheye.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle
  SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='fisheye'>
  <service name='application/fisheye' type='service' version='1'>
    <create_default_instance enabled='false' />

    <dependency
      name='multi-user-server'
      grouping='optional_all'
      type='service'
      restart_on='none'>
        <service_fmri value='svc:/milestone/multi-user-server' />
    </dependency>

    <exec_method type='method' name='start'
      exec='/opt/wmf/smf/method-fisheye %m'
      timeout_seconds='60'>
      <method_context>
        <method_credential user='fisheye' />
      </method_context>
    </exec_method>

    <exec_method type='method' name='stop' exec='/opt/wmf/smf/method-fisheye %m' timeout_seconds='60'>
      <method_context>
        <method_credential user='fisheye' />
      </method_context>
    </exec_method>

      <property_group name='startd' type='framework'>
	<propval name='duration' type='astring' value='transient' />
      </property_group>

    <template>
      <common_name><loctext xml:lang='C'>FishEye</loctext></common_name>
    </template>
  </service>
</service_bundle>
