# tools/smf/sshtunnel.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle
  SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='sshtunnel'>
  <service name='network/sshtunnel' type='service' version='1'>
    <create_default_instance enabled='false' />

    <dependency
      name='multi-user-server'
      grouping='optional_all'
      type='service'
      restart_on='none'>
        <service_fmri value='svc:/milestone/multi-user-server' />
    </dependency>

    <exec_method type='method' name='start' exec='/opt/wmf/smf/method-sshtunnel %m' timeout_seconds='60'>
    </exec_method>

    <exec_method type='method' name='stop' exec=':kill' timeout_seconds='60'>
    </exec_method>

      <property_group name='startd' type='framework'>
	<propval name='duration' type='astring' value='child' />
      </property_group>

         <property_group name='tunnel' type='application'>
                 <propval name='local-port' type='integer' value='0' />
                 <propval name='remote-port' type='integer' value='0' />
                 <propval name='remote-host' type='astring' value='' />
         </property_group>

    <template>
      <common_name><loctext xml:lang='C'>SSH tunnel</loctext></common_name>
    </template>
  </service>
</service_bundle>
