# rpms/SOURCES/squid.cron
# Restart squid if it crashed
* * * * *	root	/sbin/pidof /usr/sbin/squid > /dev/null || ([ -e /var/lock/subsys/squid ] && /sbin/service squid start)
# Rotate logs every 10 minutes
*/10 * * * *	root	/usr/sbin/logrotate -f /etc/logrotate.d/squid >/dev/null 2>&1
